# gateway/mesh-portatil-gateway.service
[Unit]
Description=Gateway Mesh Portatil (chat + @claude)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pfac
WorkingDirectory=/home/pfac/mesh-portatil/gateway
EnvironmentFile=/home/pfac/mesh-portatil/gateway/.env
ExecStart=/home/pfac/mesh-portatil/gateway/.venv/bin/python gateway.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
